# ilegra/ilegra-jovens-talentos/3-devops/tema13/scripts/calculator.service
[Unit]
Description=Golang Calculator

[Service]
User=ec2-user
ExecStart=/home/ec2-user/scripts/run.sh

[Install]
WantedBy=default.target